# php-fpm/supervisord.conf
[supervisord]
nodaemon=true

[program:php7.0-fpm]
command=/usr/sbin/php-fpm7.2 -c /etc/php/7.2/fpm
autostart=true
user=root
